# bin/cron.daily
#!/usr/bin/env php
<?php

/* cron.daily - maintenance tasks that should ideally be executed on a
 * daily basis */

require __DIR__.'/../inc/cron.php';

\Osmium\Cron\run(
	'cron.daily',
	function() {
		passthru('./bin/prune_database');
		passthru('make update-eveapis');
		passthru('make update-affiliations');
		passthru('./bin/reindex_loadouts');
		passthru('./bin/make_sitemap');
	},
	4 * 3600,
	8 * 3600
);
