# bin/cron.hourly
#!/usr/bin/env php
<?php

/* cron.hourly - maintenance tasks that should ideally be executed on
 * an hourly basis */

require __DIR__.'/../inc/cron.php';

\Osmium\Cron\run(
	'cron.hourly',
	function() {
		passthru('./bin/update_lscores');
		passthru('./bin/prune_cache');
		passthru('./bin/cache_top_kills');
	},
	20 * 60,
	30 * 60
);
