# ext/osmium-eve-kill-stomp.supervisord.conf
[program:osmium-eve-kill-stomp]
command=/srv/http/osmium/bin/eve_kill_stomp_client
autostart=true
autorestart=unexpected
exitcodes=0
startretries=10
user=osmium
